# templates/default/chef-init.erb
#!/bin/bash
#
# Startup script for the service <%= @srv_name %> (Chef Log Monitoring)
#
# $Id: Generated by chef for node: <%= node['hostname'] %>
# ----------------------------------------------------------------
# NOTE: This file is controlled by chef templates!
# Do not edit or change this file but change the following:
# template file : <%= @template_file %>
# recipe file   : <%= @recipe_file %>
# ----------------------------------------------------------------
# Schuberg Philis 2014
# ----------------------------------------------------------------
# chkconfig: - 98 02
#
### BEGIN INIT INFO
# Provides: <%= @srv_name %>
# Required-Start: $local_fs $network $remote_fs
# Required-Stop: $local_fs $network $remote_fs
# Should-Start: $named $time
# Should-Stop: $named $time
# Short-Description: Startup script for the chef monitoring service
# Description: Monitoring frontend / backend client script for manage the service.
### END INIT INFO

exec="<%= @bin_dir %>/<%= @srv_name %>"

start() {
    $exec start -- -C "<%= @install_dir %>/config.rb"
    return $?
}

stop() {
    $exec stop -- -C "<%= @install_dir %>/config.rb"
    return $?
}

restart () {
    $exec restart -- -C "<%= @install_dir %>/config.rb"
    return $?
}

reload() {
    $exec reload -- -C "<%= @install_dir %>/config.rb"
    return $?
}

populate() {
    $exec start -- -C "<%= @install_dir %>/config.rb" -O $1
    return $?
}

force_reload() {
    restart
    return $?
}

rh_status() {
    $exec status -- -C "<%= @install_dir %>/config.rb"
    return $?

}

rh_status_q() {
    $exec status -- -C "<%= @install_dir %>/config.rb"
    return $?
}

case "$1" in
    start)
        $1
        ;;
    stop)
        $1
        ;;
    restart)
        $1
        ;;
    reload)
        $1
        ;;
    populate)
        stop
        $1 $2
        start
        ;;
    force-reload)
        force_reload
        ;;
    status)
        rh_status
        ;;
    condrestart|try-restart)
        rh_status_q
        restart
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|populate|condrestart|try-restart|reload|force-reload}"
        exit 2
esac
exit $?
